# semaphore/rootfs/etc/s6-overlay/s6-rc.d/init-semaphore/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Semaphore
# Semaphore init script, runs before any other service
# ==============================================================================

bashio::log.info "Configuring Semaphore..."

# Ensure configuration data directory exists
if ! bashio::fs.directory_exists '/config/data'; then
  mkdir -p /config/data || bashio::exit.nok "Failed to create semaphore configuration directory"
fi

if ! bashio::fs.directory_exists '/data/tmp'; then
  mkdir -p /data/tmp || bashio::exit.nok "Failed to create semaphore tmp directory"
fi

bashio::var.json \
  public_url "$(bashio::config 'public_url')" \
  db_type "$(bashio::config 'db_type')" \
  db_name "$(bashio::config 'db_name')" \
  db_host "$(bashio::config 'db_host')" \
  db_port "^$(bashio::config 'db_port')" \
  db_user "$(bashio::config 'db_user')" \
  db_password "$(bashio::config 'db_password')" \
  | tempio \
      -template /etc/semaphore/config.gtpl \
      -out /config/config.json

if [ "$(bashio::config 'log_level')" == "debug" ]; then
  bashio::log.info "$(cat /config/config.json)"
fi

# Create default user
users=$(semaphore --config=/config/config.json user list | grep -E '^[a-zA-Z0-9_-]+$')

# Check if the output is empty
if [ -z "$users" ]; then
  echo "No users found. Creating a new user..."

  semaphore --config=/config/config.json user add --admin --login admin --name Admin --email admin@example.com --password admin
else
  echo "Existing users:"
  echo "$users"
fi

