# semaphore/rootfs/etc/s6-overlay/s6-rc.d/semaphore/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: NSemaphore8N
# Runs Semaphore script
# ==============================================================================

bashio::log.info "Starting Semaphore..."

exec /user/local/bin/semaphore service --config=/config/semaphore/config.json
