# semaphore/rootfs/etc/s6-overlay/s6-rc.d/init-semaphore-user/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Semaphore
# Semaphore init script, runs before any other service
# ==============================================================================

bashio::net.wait_for 43002

bashio::log.info "Configuring Semaphore Default User ..."



