# semaphore/rootfs/etc/s6-overlay/s6-rc.d/init-semaphore-user/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Semaphore
# Semaphore init script, runs before any other service
# ==============================================================================

bashio::net.wait_for 43002

bashio::log.info "Configuring Semaphore Default User ..."
declare users

# Create default user
users=$(semaphore --config=/data/config.json user list | grep -E '^[a-zA-Z0-9_-]+$')

# Check if the output is empty
if [[ -z "$users" ]]; then
  bashio::log.info "No users found. Creating a new user..."

  semaphore --config=/data/config.json user add --admin --login admin --name Admin --email admin@example.com --password admin
else
  bashio::log.info "Existing users:"
  bashio::log.info "$users"
fi

