# semaphore/rootfs/etc/s6-overlay/s6-rc.d/init-semaphore-user/up
/etc/s6-overlay/s6-rc.d/init-semaphore-user/run
