# semaphore/rootfs/etc/s6-overlay/s6-rc.d/init-python/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Semaphore
# Configures Python packages
# ==============================================================================

# Install Python packages, if enabled.
if bashio::config.has_value 'python_packages'; then
    for package in $(bashio::config 'python_packages' | jq -r '.[]'); do
        pip install "${package}"
    done
fi 
